# CONTROL/preinst
#!/usr/bin/env ash
#

systemctl stop watchdog.timer
systemctl stop watchdog.service

systemctl disable watchdog.timer
systemctl disable watchdog.service
